# .config/supervisord.conf
[program:uwsgi]
command = uwsgi --ini /srv/projects/.config/uwsgi.ini

[program:nginx]
command = nginx -g 'daemon off;'